# battery-monitor.timer
[Unit]
Description=Battery monitor

[Timer]
OnStartupSec=10
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
